# docker/root/etc/s6/gitea/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/shipyard >/dev/null
exec su-exec $USER /usr/local/bin/shipyard web
popd
